# nanomq/rootfs/etc/services.d/nanomq/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs NanoMQ
# ==============================================================================

declare -a options

options+=(start)

if bashio::debug; then
    options+=(--log_level debug)
else
    options+=(--log_level info)
fi

options+=(--http start)

bashio::log.info "Starting NanoMQ MQTT broker..."
exec nanomq "${options[@]}"
